# infra/supervisor/pennant-sse.conf
[program:pennant-sse]
process_name=%(program_name)s
command=/home/ubuntu/.nvm/versions/node/v22.22.2/bin/node /var/www/pennant-sse/server.js
directory=/var/www/pennant-sse
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=ubuntu
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/pennant-sse.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=NODE_ENV="production",PENNANT_SSE_PORT="8081",REDIS_URL="redis://127.0.0.1:6379",PENNANT_API_BASE="http://127.0.0.1",PENNANT_PUBSUB_PREFIX="flag-updates"
